# script/supervisor.conf
[program:air]
directory=/src/app
command=air
process_name=%(program_name)s_%(process_num)02d
numprocs=1
redirect_stderr=true
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor.log